# protocol_pulse.service
[Unit]
Description=Protocol Pulse Flask Server (Gunicorn)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ultron/protocol_pulse
Environment=PYTHONUNBUFFERED=1
Environment=FLASK_ENV=production
Environment=ENABLE_SELF_CHECK_BYPASS=true
EnvironmentFile=-/home/ultron/protocol_pulse/.env
ExecStart=/home/ultron/protocol_pulse/venv/bin/gunicorn -w 1 -k gthread --threads 8 -b 0.0.0.0:5000 app:app
Restart=always
RestartSec=5
StandardOutput=append:/home/ultron/protocol_pulse/logs/app.log
StandardError=append:/home/ultron/protocol_pulse/logs/app.log

[Install]
WantedBy=default.target
